# scripts/sattrack-serve.service
# SatTrack web dashboard — real-time globe + telemetry (port 8082).
#
# Install:
#   sudo bash scripts/install_serve_service.sh

[Unit]
Description=SatTrack web dashboard (FastAPI + WebSocket on :8082)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/jellyboy/MCP-Kali/SatTrack
ExecStart=/usr/bin/python3 /home/jellyboy/MCP-Kali/SatTrack/run.py serve --host 0.0.0.0 --port 8082
Restart=on-failure
RestartSec=5
Environment=SATTRACK_DATA=/var/lib/sattracker

[Install]
WantedBy=multi-user.target
